# root/var/service/sshd/run
#!/bin/bash

export HOME=/root

mkdir -p /etc/sshd

ssh-keygen -A >/dev/null 2>&1

gpg --import /etc/keys/*.pub

for fingerprint in ${ADMIN_FINGERPRINTS//,/ }; do
	add-key admin "${fingerprint}"
	add-key user "${fingerprint}"
done
for fingerprint in ${USER_FINGERPRINTS//,/ }; do
	add-key user "${fingerprint}"
done

exec /usr/sbin/sshd \
	-f /var/service/sshd/sshd_config \
	-D \
	-e \
	-o LogLevel=VERBOSE \
	2>&1
